# backend/ops/cron/laravel-scheduler.cron
* * * * * cd /var/www/hrm-nawa/backend && php artisan schedule:run >> /dev/null 2>&1
