# backend/ops/supervisor/hrm-nawa-workers.conf
[program:hrm-nawa-queue-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/hrm-nawa/backend/artisan queue:work --queue=default --sleep=3 --tries=3 --max-time=3600 --timeout=120
directory=/var/www/hrm-nawa/backend
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=deploy
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/hrm-nawa-queue-default.log
stopwaitsecs=3600

[program:hrm-nawa-queue-ai-alerts]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/hrm-nawa/backend/artisan queue:work --queue=ai-alerts --sleep=2 --tries=5 --backoff=60 --max-time=3600 --timeout=120
directory=/var/www/hrm-nawa/backend
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=deploy
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/hrm-nawa-queue-ai-alerts.log
stopwaitsecs=3600
